# systemd/user/backup.service
[Unit]
Description=Periodical backup
Wants=daily.timer

[Service]
ExecStart=/home/maca/bin/backup

[Install]
WantedBy=daily-timer.target
