# systemd/user/daily.timer
[Unit]
Description=Daily Timer

[Timer]
OnBootSec=5min
OnCalendar=daily
Unit=daily-timer.target

[Install]
WantedBy=basic.target
